# config/systemd/user/save-notes.service
[Unit]
Description=Push notes to github

[Service]
Type=oneshot
ExecStart=%h/.config/gevious-scripts/push-notes
