# config/systemd/user/wallpaper.service
[Unit]
Description=Fetch Unsplash robots wallpaper and set in Sway

[Service]
ExecStart=%h/.config/gevious-scripts/get-wallpaper
